# telegram-bot.service
[Unit]
Description=Telegram Quiz Bot Service (Modullashtirilgan)
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/quiz-bot
Environment="PYTHONPATH=/home/azureuser/quiz-bot"
# Environment variables .env fayldan o'qiladi yoki quyidagilarni ishlatish mumkin:
# BOT_TOKEN .env fayldan o'qiladi
Environment="ADMIN_USER_IDS=7537966029"
Environment="DEEPSEEK_API_KEY=sk-2ef0592830e84305825f6344f408e226"
# Yangi cheklovlar (ixtiyoriy, default qiymatlar ishlatiladi)
# Environment="MIN_QUESTIONS_REQUIRED=1"
# Environment="MAX_QUESTIONS_PER_QUIZ=200"
# Environment="MAX_SEQUENTIAL_QUESTIONS=150"
# Environment="MAX_AI_SECONDS=300"
# Environment="MAX_TEXT_CHARS_FOR_AI=50000"
ExecStart=/home/azureuser/quiz-bot/venv/bin/python3 -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Memory va CPU cheklovlari (ixtiyoriy)
# MemoryLimit=512M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
